# deploy/circuitikz-render.service
[Unit]
Description=CircuitikZ CAD render server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/circuitikz-cad/current
ExecStart=/usr/bin/node /var/www/circuitikz-cad/current/render-server/server.mjs
Restart=always
RestartSec=2
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
